# crontab
MAILTO=doloopwhile@gmail.com
00,15,30,45 * * * * /var/bot/bot.sh
